# solarnode-deploy/equinox/base/bin/solarssh-cleaner.timer
[Unit]
Description=Close discarded SolarSSH sessions hourly.

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
